# deploy/discord-codex-bot.service
[Unit]
Description=Discord Codex Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/user/projects/Nanihold_OS
EnvironmentFile=/home/user/projects/Nanihold_OS/.env
ExecStart=/home/user/projects/Nanihold_OS/.venv/bin/python /home/user/projects/Nanihold_OS/bot/discord_codex_bot.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
